# jvb/rootfs/etc/s6-overlay/s6-rc.d/jvb/run
#!/command/execlineb -P
s6-setuidgid jvb
/etc/s6-overlay/scripts/jvb
