# prosody/rootfs/etc/s6-overlay/s6-rc.d/saslauthd/run
#!/command/execlineb -P
s6-setuidgid root
with-contenv
saslauthd -a ldap -O /etc/saslauthd.conf -c -m /var/run/saslauthd -n 5 -d
