# prosody/rootfs/etc/s6-overlay/s6-rc.d/prosody/run
#!/command/execlineb -P
s6-setuidgid prosody
with-contenv
prosody --config /config/prosody.cfg.lua -F
