# jibri/rootfs/etc/s6-overlay/s6-rc.d/icewm/run
#!/command/execlineb -P
s6-setuidgid jibri
with-contenv
/usr/bin/icewm-session
