# jigasi/rootfs/etc/s6-overlay/s6-rc.d/jigasi/run
#!/command/execlineb -P
s6-setuidgid jigasi
/etc/s6-overlay/scripts/jigasi
