# jicofo/rootfs/etc/s6-overlay/s6-rc.d/jicofo/run
#!/command/execlineb -P
s6-setuidgid jicofo
/etc/s6-overlay/scripts/jicofo
