# web/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb
/etc/s6-overlay/scripts/cron
