# jibri/rootfs/etc/s6-overlay/s6-rc.d/xorg/run
#!/command/execlineb -P
s6-setuidgid jibri
/etc/s6-overlay/scripts/xorg
